# rootfs/etc/s6-overlay/s6-rc.d/30-janus-config/up
/etc/s6-overlay/scripts/30-janus-config.sh