# rootfs/etc/s6-overlay/s6-rc.d/janus-gateway/finish
#!/command/with-contenv bash
. "/usr/local/bin/logger"

echo "Bringing container down NOW!"

/run/s6/basedir/bin/halt